# deploy/roles/svr_dev/templates/supervisord.conf
[program:{{ program_name }}]
command={{ app_install_dir }}/env/bin/uwsgi --die-on-term --ini {{ app_install_dir }}/conf/uwsgi.ini
user=www
group=www
